# device/testing/sheng-mipps/90-xiaomi-mipps-auth.rules
# Trigger Xiaomi MiPPS authentication after USB-C partner attach or
# USB power becoming available through an already-attached dock.
# The long-running work is delegated to systemd; udev only schedules the oneshot service.

ACTION=="add", SUBSYSTEM=="typec", KERNEL=="port*-partner", TAG+="systemd", ENV{SYSTEMD_WANTS}+="xiaomi-mipps-auth.service"
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="qcom-battmgr-usb", ENV{POWER_SUPPLY_ONLINE}=="1", RUN+="/usr/libexec/xiaomi-mipps-auth --usb-online-trigger"
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="qcom-battmgr-usb", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/usr/libexec/xiaomi-mipps-auth --clear-notification"
