# device/testing/sheng-mipps/xiaomi-mipps-auth.service
[Unit]
Description=Xiaomi MiPPS/PPS charger authentication
ConditionPathExistsGlob=/sys/devices/platform/pmic-glink/*/xiaomi/request_vdm_cmd

[Service]
Type=oneshot
ExecStart=/usr/bin/flock -n -E 0 /run/xiaomi-mipps-auth.lock /usr/libexec/xiaomi-mipps-auth --timeout 4
TimeoutStartSec=20
