# device/community/device-nokia-n900/device-nokia-n900.initd
#!/sbin/openrc-run

command="/usr/sbin/alsactl"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="restore -f /var/lib/alsa/asound.state.speakers"
