# device/community/device-pine64-pinebookpro/wlan-powersave-off.initd
#!/sbin/openrc-run

command="/usr/sbin/iw"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="dev wlan0 set power_save off"
