# main/postmarketos-usb-moded/rootfs-etc-init.d-usb-moded-tethering-mode.initd
#!/sbin/openrc-run

description="USB Moded Tethering Mode"
command="usb-moded-tethering-mode"
command_args="up"

depends() {
    need net
}

stop() {
    usb-moded-tethering-mode down
}
