# main/postmarketos-usb-moded/rootfs-etc-init.d-usb-moded-mtp-mode.initd
#!/sbin/openrc-run

description="USB Moded MTP Mode"
command="usb-moded-mtp-mode"
command_args="up"

stop() {
    usb-moded-mtp-mode down
}
